# docker/spring/etc/supervisor/conf.d/sshd.conf
[program:sshd]
command=/usr/bin/sshd -D -4 -p 10022
autostart=true
autorestart=true
user=root
