# deploy/roles/mint-v2.1/templates/v2/scheduler.timer.j2
[Unit]
Description=Runs Mint {{ app_domain }} scheduler every minute

[Timer]
OnBootSec=15minutes
OnUnitActiveSec=1minute
Unit=scheduler.{{ app_domain }}.service

[Install]
WantedBy=timers.target
