# deploy/roles/mint-v2.0/templates/scheduler.service.j2
[Unit]
Description=Mint scheduler for {{ app_domain }}

[Service]
Type=simple
ExecStart=/usr/bin/php artisan schedule:run
WorkingDirectory={{ app_deploy_target }}
User=www-data
Group=www-data
